# dist/rpm/SOURCES/tapir-pop.service
[Unit]
Description=DNSTAPIR EDGE Policy Processor
Wants=unbound.service
Before=unbound.service

[Service]
Type=simple
User=tapir-edge
Group=tapir-edge
ExecStart=/usr/bin/tapir-pop

[Install]
WantedBy=multi-user.target
